# units/postgresql@.service.d/module-pg_trim.conf
# Adds the pg_trgm module.
# Note that this will override any other modules specified in the DB_EXTENSION variable.
[Service]
Environment='DOCKER_ENV_DB_EXTENSION=DB_EXTENSION=pg_trgm'
