# units/postgresql@.service.d/publish-5432.conf
[Service]
# Expose the PostgreSQL port (5432)
Environment='DOCKER_PUBLISH_5432=5432:5432'
